# edge_node_service/edge-node.service
# systemd 服务单元：用于在 Linux 板卡启动时自动拉起 edge_node_service。
# 工作目录、可执行文件路径和重启策略需与实际部署目录保持一致。
[Unit]
Description=FPGA Edge Vision Node Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/pi/edge_node_service
ExecStart=/home/pi/edge_node_service/build/edge_node_service /home/pi/edge_node_service/config/config.json
Restart=always
RestartSec=2
User=pi

[Install]
WantedBy=multi-user.target
